# infra/systemd/nidalheim-api-auth-staging.service
[Unit]
Description=Nidalheim api-auth (staging, host process)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ubuntu/Nidalheim-staging/services/api-auth
Environment=NIDALHEIM_ROOT=/home/ubuntu/Nidalheim-staging
ExecStart=/home/ubuntu/Nidalheim-staging/services/api-auth/start-staging.sh
Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5
StandardOutput=journal
StandardError=journal
KillSignal=SIGTERM
TimeoutStopSec=15

[Install]
WantedBy=default.target
